# gateway/smsBEST.conf
# this file goes in /etc/supervisor/conf.d
[program:smsBEST]
command=/srv/smsBEST/bin/run_gateway.sh
directory=/srv/smsBEST/var
autorestart=true
